# ipwatchd/com.quest.rc.ipwatchd.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.
com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>OnDemand</key>
        <false/>

        <key>Label</key>
        <string>com.quest.rc.ipwatchd</string>

        <key>ProgramArguments</key>
        <array>
	    <string>@sbindir@/ipwatchd</string>
	    <string>@sbindir@/dnsupdate</string>
        </array>

        <key>ServiceDescription</key>
        <string>Synchronizes DNS with IP address changes</string>

	<key>Debug</key>
	<false/>

	<key>SHAuthorizationRight</key>
	<string>system.preferences</string>

	<key>ServiceIPC</key>
	<true/>
</dict>
</plist>

